# config/systemd/user/cpu-quota.slice
[Unit]
Description= A slice used to scope applications which run away with a laptop's battery

[Slice]
CPUAccounting=true
CPUQuota=300%

