# config/systemd/user/dropbox.service
[Unit]
Description=Launches DropBox Flatpack 

[Service]
Type=exec
ExecStart=flatpak run com.dropbox.Client

[Install]
WantedBy=default.target

